# 99-video.rules
SUBSYSTEM=="video4linux", ACTION=="add", MODE:="0777", RUN+="/bin/su pi -c '/etc/udev/rules.d/webcam.sh %k'"
